# systemd/.config/systemd/user/restic-backup.service
[Unit]
Description=Daily Restic Backup

[Service]
Type=oneshot
ExecStart=%h/bin/sh/backup-restic-to-data